# provision/ansible/roles/bridge/templates/bridge.service.j2
# Bridge Service systemd Unit
# Generated by Ansible - Do not edit manually
# {{ ansible_managed }}

[Unit]
Description=Busibox Bridge Service (Multi-channel Communication)
Documentation=https://github.com/your-org/busibox
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User={{ bridge_user }}
Group={{ bridge_group }}
WorkingDirectory={{ bridge_dir }}
EnvironmentFile={{ bridge_dir }}/.env

ExecStart={{ bridge_venv }}/bin/python -m app.main

Restart=always
RestartSec=10
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ bridge_service_name }}

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ bridge_dir }}
PrivateTmp=true

[Install]
WantedBy=multi-user.target
