# provision/ansible/roles/app_deployer/templates/app-systemd.service.j2
# {{ ansible_managed }}
# Systemd service for {{ app.name }}
# Generated by app_deployer role

[Unit]
Description={{ app.name }} Application
After=network.target
{% if app.requires is defined %}
{% for req in app.requires %}
Requires={{ req }}.service
After={{ req }}.service
{% endfor %}
{% endif %}

[Service]
Type=simple
User=root
Group=root
WorkingDirectory={{ app.deploy_path }}

# Load environment from .env file
EnvironmentFile=-{{ app.deploy_path }}/.env

# Environment variables from apps.yml
{% if app.env is defined %}
{% for key, value in app.env.items() %}
Environment="{{ key }}={{ value }}"
{% endfor %}
{% endif %}

# Port configuration
Environment="PORT={{ app.port }}"
Environment="NODE_ENV=production"
Environment="HOSTNAME=0.0.0.0"
# Cap V8 heap to prevent unbounded memory growth
Environment="NODE_OPTIONS=--max-old-space-size=512"

# Start command
# In a monorepo standalone build, server.js may be nested under the app's
# relative path (e.g. .next/standalone/apps/portal/server.js).
# We find it dynamically to support both monorepo and standalone layouts.
ExecStart=/bin/bash -c '\
  SERVER_JS=$(find .next/standalone -name server.js -maxdepth 4 2>/dev/null | head -1); \
  if [ -n "$SERVER_JS" ]; then exec node "$SERVER_JS"; \
  elif [ -f .next/standalone/server.js ]; then exec node .next/standalone/server.js; \
  else exec npm start; fi'

# Restart configuration
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# Logging - use journal for consistent timestamps
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ app.name }}

# Security hardening (optional, can be enabled per-app)
# NoNewPrivileges=true
# ProtectSystem=strict
# ProtectHome=true
# PrivateTmp=true

[Install]
WantedBy=multi-user.target

