# provision/ansible/roles/data_api/templates/data-api.service.j2
[Unit]
Description=Busibox Data API Server
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target
Requires=postgresql.service redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/srv/data-api/src
EnvironmentFile=/srv/data-api/.env
ExecStart=/srv/data-api/venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8002
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

